# .config/systemd/user/xdg-desktop-portal-gnome.service
[Unit]
ConditionEnvironment=DESKTOP_SESSION!=gnome
ConditionEnvironment=DESKTOP_SESSION!=gdm
[Service]
ExecStart=
